# insa-iot-platform/scripts/timescaledb-backup.timer
[Unit]
Description=Daily TimescaleDB Backup Timer
Documentation=file:///home/wil/insa-iot-platform/docs/BACKUP_AUTOMATION.md
Requires=timescaledb-backup.service

[Timer]
# Run daily at 2 AM
OnCalendar=*-*-* 02:00:00

# Run on boot if missed (e.g., server was down)
Persistent=true

# Randomize start time by up to 15 minutes to avoid thundering herd
RandomizedDelaySec=15min

# Log timer events
Unit=timescaledb-backup.service

[Install]
WantedBy=timers.target
